# examples/full-stack/generated/quadlet-network-frontend
#!/sbin/openrc-run
# Generated by quadlet-openrc 0.3.0 from frontend.network
# source-hash: f7c4743b35a725a8
# DO NOT EDIT: change the Quadlet source and regenerate.

description='Frontend application network'
depend() {
    need podman
}

start() { /usr/bin/podman network exists 'frontend' || /usr/bin/podman 'network' 'create' '--driver' 'bridge' 'frontend'; }
stop() { :; }
status() { /usr/bin/podman network exists 'frontend'; }
