# examples/full-stack/generated/quadlet-volume-webdata
#!/sbin/openrc-run
# Generated by quadlet-openrc 0.3.0 from webdata.volume
# source-hash: fb4e06244b86b962
# DO NOT EDIT: change the Quadlet source and regenerate.

description='Persistent web content'
depend() {
    need podman
}

start() { /usr/bin/podman volume exists 'webdata' || /usr/bin/podman 'volume' 'create' 'webdata'; }
stop() { :; }
status() { /usr/bin/podman volume exists 'webdata'; }
