# examples/full-stack/generated/quadlet-image-nginx
#!/sbin/openrc-run
# Generated by quadlet-openrc 0.3.0 from nginx.image
# source-hash: 694a600841d5ad5b
# DO NOT EDIT: change the Quadlet source and regenerate.

description='Pinned nginx image'
depend() {
    need podman
}

start() { /usr/bin/podman image exists 'localhost/systemd-nginx' || { /usr/bin/podman 'pull' 'docker.io/library/nginx:1.29-alpine' && /usr/bin/podman tag 'docker.io/library/nginx:1.29-alpine' 'localhost/systemd-nginx'; }; }
stop() { :; }
status() { /usr/bin/podman image exists 'localhost/systemd-nginx'; }
