# examples/full-stack/generated/web
#!/sbin/openrc-run
# Generated by quadlet-openrc 0.3.0 from web.container
# source-hash: 1e5289d09ca52c2e
# DO NOT EDIT: change the Quadlet source and regenerate.

description='Hardened nginx example'
name='web'
supervisor=supervise-daemon
command='/bin/sh'
command_args='-c exec /usr/bin/podman '\''run'\'' '\''--rm'\'' '\''--replace'\'' '\''--name'\'' '\''web'\'' '\''--read-only'\'' '\''--volume'\'' '\''webdata:/usr/share/nginx/html:rw'\'' '\''--publish'\'' '\''8080:80'\'' '\''--cap-drop'\'' '\''ALL'\'' '\''--network'\'' '\''frontend'\'' '\''--tmpfs'\'' '\''/tmp:rw,noexec,nosuid,size=64m'\'' '\''--tmpfs'\'' '\''/run:rw,noexec,nosuid,size=16m'\'' '\''--userns'\'' '\''auto'\'' '\''--pids-limit'\'' '\''256'\'' '\''--memory'\'' '\''512m'\'' '\''--cpus'\'' '\''1'\'' '\''--health-cmd'\'' '\''wget -q -O /dev/null http://127.0.0.1/'\'' '\''--health-interval'\'' '\''30s'\'' '\''--health-timeout'\'' '\''5s'\'' '\''--health-retries'\'' '\''3'\'' '\''--security-opt'\'' '\''no-new-privileges'\'' '\''localhost/systemd-nginx'\'''
output_log=/var/log/quadlet/web.log
error_log=/var/log/quadlet/web.err
respawn_delay=3
respawn_max=0

depend() {
    need podman
    need quadlet-image-nginx quadlet-network-frontend quadlet-volume-webdata
    after quadlet-image-nginx quadlet-network-frontend quadlet-volume-webdata
}

start_pre() { checkpath -d -m 0755 /var/log/quadlet; }
stop_post() { /usr/bin/podman stop -t '15' 'web' >/dev/null 2>&1 || true; }
